# server-backup.timer
[Unit]
Description=Run backup script weekly

[Timer]
OnCalendar=Mon *-*-* 8:00:00
Persistent=true

[Install]
WantedBy=timers.target